# deployment/systemd/ai-scheduler.service
[Unit]
Description=AI Scheduler (Embedding Indexes)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/vinay/security-ops-platform
ExecStart=/home/vinay/security-ops-platform/.venv/bin/python3 -u src/ai_scheduler.py
Restart=on-failure
KillMode=mixed
TimeoutStopSec=10
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
