# deployment/systemd/ir-claude-router-shim.service
[Unit]
Description=IR Claude Router Shim (front door for ir-claude-router; adds /v1/models discovery, friendly-id rewrite, bearer auth)
After=network.target ir-claude-router.service
Wants=ir-claude-router.service

[Service]
Type=simple
WorkingDirectory=/home/vinay/security-ops-platform
EnvironmentFile=/home/vinay/security-ops-platform/data/transient/.env
ExecStart=/home/vinay/security-ops-platform/.venv/bin/python deployment/claude_router_shim.py
Restart=on-failure
KillMode=mixed
TimeoutStopSec=10
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/vinay/security-ops-platform/data/transient/logs/claude_router_shim.log
StandardError=append:/home/vinay/security-ops-platform/data/transient/logs/claude_router_shim.log

[Install]
WantedBy=default.target
