# deployment/systemd/ir-logdy-money-ball.service
[Unit]
Description=Logdy - MoneyBall Logs
After=network-online.target ir-money-ball.service
Wants=network-online.target

[Service]
Type=simple
User=vinay
Group=vinay
WorkingDirectory=/home/vinay/pub/IR
ExecStart=/bin/bash -c 'tail -f /home/vinay/pub/IR/logs/money_ball.log | /home/vinay/bin/logdy stdin --ui-ip 0.0.0.0 -p 8034 --ui-pass=$(cat /home/vinay/pub/IR/.logdy_password 2>/dev/null || echo "metcirt")'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ir-logdy-money-ball

[Install]
WantedBy=multi-user.target
