# deployment/systemd/ir-peer-ping.service
[Unit]
Description=IR Peer Ping - Keep bot NAT paths active
After=network.target

[Service]
Type=oneshot
User=user
Group=user
WorkingDirectory=/home/user/pub/IR
Environment="PYTHONPATH=/home/user/pub/IR"
ExecStart=/home/user/pub/IR/.venv/bin/python /home/user/pub/IR/src/peer_ping_keepalive.py
StandardOutput=append:/home/user/pub/IR/logs/peer_ping.log
StandardError=append:/home/user/pub/IR/logs/peer_ping.log

# Resource limits
MemoryMax=100M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
