# deployment/systemd/ir-peer-ping.service
[Unit]
Description=IR the peer-status service Keepalive
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/vinay/security-ops-platform
ExecStart=/home/vinay/security-ops-platform/.venv/bin/python3 src/peer_ping_keepalive.py
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/vinay/security-ops-platform

[Install]
WantedBy=default.target
